# deploy/chatagent.service
[Unit]
Description=ChatAgentCore Service
After=network.target

[Service]
# 指向打包后的可执行文件路径（用户安装时需修改此处）
ExecStart=/opt/chatagent/chatagent-service

# 工作目录（建议与可执行文件同级，以便自动生成 config）
WorkingDirectory=/opt/chatagent

# 总是自动重启
Restart=always
RestartSec=5

# 日志输出到 syslog
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=chatagent

[Install]
WantedBy=multi-user.target
